# My_Git_Key/60-my-git-usb-key.rules
ACTION=="add" KERNEL=="sd[a-z]*",  RUN+="/usr/bin/my_udev.bash add %k"
ACTION=="remove" KERNEL=="sd[a-z]*",  RUN+="/usr/bin/my_udev.bash remove %k"

